# data/rc.local
#!/bin/bash

/data/home/root/wg-rc.local.sh &
